# gpt_realtime_assist/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
set -e
/opt/app/bootstrap.sh
exec python3 /opt/app/main.py
